# sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tigervnc]
# Launch the Xtigervnc server directly (binary name on Debian bookworm)
command=/usr/bin/Xtigervnc :1 -geometry 1280x800 -depth 24 -SecurityTypes None -localhost no -alwaysshared
autorestart=false
priority=10
stdout_logfile=/var/log/vnc.log
stderr_logfile=/var/log/vnc.log

[program:openbox]
command=/bin/bash -c "sleep 3 && DISPLAY=:1 /usr/bin/openbox --config-file /root/.config/openbox/rc.xml"
autorestart=true
priority=20
stdout_logfile=/var/log/openbox.log
stderr_logfile=/var/log/openbox.log

[program:novnc]
# websockify proxies WebSocket → TCP(VNC port 5901); noVNC JS client connects via 6080
command=/bin/bash -c "sleep 4 && /usr/bin/websockify --web /opt/novnc 6080 localhost:5901"
autorestart=true
priority=30
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.log

[program:firefox]
command=/bin/bash -c "sleep 8 && DISPLAY=:1 /usr/bin/firefox-esr --new-window file:///sandbox/email/index.html --no-sandbox"
autorestart=false
priority=40
stdout_logfile=/var/log/firefox.log
stderr_logfile=/var/log/firefox.log

# Open Thunar file manager in Downloads so analyst can access downloaded attachments
[program:thunar]
command=/bin/bash -c "sleep 12 && DISPLAY=:1 /usr/bin/thunar /root/Downloads"
autorestart=false
priority=45
stdout_logfile=/var/log/thunar.log
stderr_logfile=/var/log/thunar.log

[program:ready-signal]
command=/bin/bash -c "sleep 14 && echo READY > /sandbox/ready && echo 'Sandbox ready'"
autorestart=false
priority=50
stdout_logfile=/var/log/ready.log
stderr_logfile=/var/log/ready.log
